# scripts/init_scripts/hidra@.service.d/transfer-client.conf
[Unit]
Wants=hidra-transfer-client@%i.service
